# docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor/

[program:rclone]
command=/usr/bin/rclone rcd --rc-addr=%(ENV_RCLONE_RC_ADDR)s --rc-user=%(ENV_RCLONE_RC_USER)s --rc-pass=%(ENV_RCLONE_RC_PASS)s --rc-allow-origin="*" --config=%(ENV_RCLONE_CONFIG)s
autostart=true
autorestart=true
startretries=3
redirect_stderr=true
stdout_logfile=/var/log/supervisor/rclone.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=3
environment=RCLONE_RC_ADDR=%(ENV_RCLONE_RC_ADDR)s,RCLONE_RC_USER=%(ENV_RCLONE_RC_USER)s,RCLONE_RC_PASS=%(ENV_RCLONE_RC_PASS)s,RCLONE_CONFIG=%(ENV_RCLONE_CONFIG)s

[program:backend]
command=/app/backend --bind=0.0.0.0:%(ENV_BACKEND_PORT)s --rclone-url=http://localhost:5572 --username=%(ENV_RCLONE_RC_USER)s --password=%(ENV_RCLONE_RC_PASS)s
autostart=true
autorestart=true
startretries=3
redirect_stderr=true
stdout_logfile=/var/log/supervisor/backend.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=3
environment=BACKEND_PORT=%(ENV_BACKEND_PORT)s,RCLONE_RC_USER=%(ENV_RCLONE_RC_USER)s,RCLONE_RC_PASS=%(ENV_RCLONE_RC_PASS)s
startsecs=5

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
startretries=3
redirect_stderr=true
stdout_logfile=/var/log/supervisor/nginx.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=3
startsecs=3
